# debian/cron.d
0 */4	* * *	root	run-parts /etc/rebuildd/maintenance
